# systemd/blinkmon-stop.service
[Unit]
Description="Stop blinkmon service"

[Service]
Type=oneshot
ExecStart=/bin/echo ''
